# etc/supervisor/conf.d/nginx.conf
[program:nginx]
command = nginx -g "daemon off;"
autostart = true
autorestart = true
stopsignal = TERM
priority = 300
